# services/sample-worker/run
#!/usr/bin/env bash

exec 2>&1

export AMQP_QUEUE=http-worker
exec  php /var/www/entrypoint/worker.php